# systemd/fwdesktop.service
[Unit]
Description=fwdesktop

[Service]
ExecStart=/usr/bin/fwdesktop

