# systemd/fwdesktopTimer.timer
[Unit]
Description=fwdesktopTimer

[Timer]
OnBootSec=1s
OnUnitActiveSec=15m
Unit=fwdesktop.service

[Install]
WantedBy=multi-user.target

