# systemd/fydesktop.service
[Unit]
Description=fydesktop

[Service]
ExecStart=/usr/bin/fydesktop

