# systemd/fydesktopTimer.timer
[Unit]
Description=fydesktopTimer

[Timer]
OnBootSec=1s
OnUnitActiveSec=15m
Unit=fydesktop.service

[Install]
WantedBy=multi-user.target

